# led-matrix-server.service
[Unit]
Description=LED Matrix System (Viewer + Web Server)
After=network.target

[Service]
Type=forking
User=root
WorkingDirectory=/home/pi/rpi-led-matrix
ExecStart=/home/pi/rpi-led-matrix/start.sh
ExecStop=/home/pi/rpi-led-matrix/stop.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
